# frameo_control/rootfs/etc/services.d/frameo/run
#!/bin/sh
# Start the Frameo ADB backend server

echo "[INFO] Starting the Frameo ADB backend server..."

# Execute the Python server directly
exec python3 /app/app.py
